# voicepipe-transcriber.service
[Unit]
Description=Voicepipe Transcriber Service
After=network.target
# Start after the main voicepipe daemon
After=voicepipe.service
Requires=voicepipe.service

[Service]
Type=simple
ExecStart=/home/pepper/.local/src/voicepipe/transcriber_daemon.py
Restart=on-failure
RestartSec=5
StandardError=journal
StandardOutput=journal

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp

# Environment
Environment="PYTHONUNBUFFERED=1"
PassEnvironment=OPENAI_API_KEY
PassEnvironment=OPENAI_API_KEY

[Install]
WantedBy=default.target